# projects/scientific-release/backend/scientific-release.service
[Unit]
Description=Scientific Release Mini-Program Backend
After=network.target mysqld.service

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/.openclaw/workspace/projects/scientific-release/backend
EnvironmentFile=/home/admin/.openclaw/workspace/projects/scientific-release/backend/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scientific-release

[Install]
WantedBy=multi-user.target
